# units/cleanup.timer
# https://www.freedesktop.org/software/systemd/man/systemd.time.html
[Unit]
Description=Run date.service daily

[Timer]
OnCalendar=daily
